# env/laravel/modules/worker/supervisor/services.d/worker/run
#!/usr/bin/execlineb -P

/usr/local/bin/php /var/www/html/artisan queue:work --tries=3
