# modules/load-balancer/supervisor/services.d/consul-template/run
#!/usr/bin/execlineb -P

/usr/bin/consul-template -config=/etc/nginx/conf.d/consul-template-config.hcl
